# etc/systemd/wsbot.service
[Unit]
Description=WebSocket Price Bot
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/4narek-info
ExecStart=/usr/bin/go run /home/ubuntu/4narek-info/main.go
Restart=always

[Install]
WantedBy=multi-user.target